# com.fire-tv-backgrounds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fire-tv-backgrounds</string>
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>__PROJECT_DIR__/server.js</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>
    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__NODE_DIR__:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
